# finance-report.service
[Unit]
Description=PIrrai - relatorio financeiro mensal por e-mail
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=rodrigor
Group=rodrigor
Environment=HOME=/home/rodrigor
WorkingDirectory=/home/rodrigor/homewatch
ExecStart=/home/rodrigor/homewatch/finance_report.py
