# routerwatch.service
[Unit]
Description=PIrrai — coleta SNMPv3 do roteador (ER605) para SQLite
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
User=rodrigor
Group=rodrigor
Environment=HOME=/home/rodrigor
WorkingDirectory=/home/rodrigor/homewatch
ExecStart=/home/rodrigor/homewatch/routerwatch.sh
